# es/files/kibana.conf
[program:kibana]
command=service kibana start
